# qca/feeds/leedarson/security/files/security.init
#!/bin/sh /etc/rc.common

# Copyright (c) Leedarson

USE_PROCD=1
START=99
STOP=13


start_service()
{
	procd_open_instance
	procd_set_param command security
	procd_set_param respawn
	procd_close_instance
}

reload_service()
{
	stop
	start
}


